# deps/xema-astermq.service
[Unit]
Description=Xema AsterMQ service

[Service]
ExecStart=/usr/bin/dotnet /var/lib/xema/astermq/AsterMQ.dll
WorkingDirectory=/var/lib/xema/astermq
User=
Group=
Restart=on-failure
SyslogIdentifier=xema-astermq
PrivateTmp=true

[Install]
WantedBy=multi-user.target